# scripts/omil.service.template
[Unit]
Description=Omil daemon

[Service]
Type=simple
ExecStart=%app_dir%/bin/omil --config %app_dir%/conf/config.yml

[Install]
WantedBy=multi-user.target
